# Objects/Threat Vectors/Scheduled tasks created with taskschlr.exe.yaml
name: Scheduled tasks created with taskschlr.exe
criticality: Medium
references:
  public:
    1: https://www.thecodebuzz.com/schedule-run-exe-console-application-windows-task-scheduler/
    2: https://attack.mitre.org/techniques/T1053/005/
    3: https://learn.microsoft.com/en-us/windows/win32/taskschd/about-the-task-scheduler
    4: https://www.windowscentral.com/how-create-automated-task-using-task-scheduler-windows-10
    5: https://www.microsoft.com/en-us/security/blog/2022/04/12/tarrask-malware-uses-scheduled-tasks-for-defense-evasion/
    6: https://redcanary.com/threat-detection-report/techniques/windows-command-shell/
  #internal:
    #a: 
  #restricted:
    #A: 
  #reports:
    #-

metadata:
  uuid: 24503678-9a1b-4af3-9837-a90bf47b7dda
  schema: tvm::2.0
  tlp: clear
  version: 2
  created: 2022-12-12
  modified: 2022-12-13
  organisation:
    uuid: 56b0a0f0-b0bc-47d9-bb46-02f80ae2065a
    name: EC DIGIT CSOC

threat:
  actors:
    - name: att&ck::G0125 #[Enterprise] HAFNIUM, Operation Exchange Marauder, Silk Typhoon
    - name: misp::4f05d6c1-3fc1-4567-91cd-dd4637cc38b5 #HAFNIUM, ATK233, G0125, Operation Exchange Marauder, Red Dev 13, Silk Typhoon
  killchain: Execution
  att&ck:
    - T1053.005
  domains:
    - Enterprise
  terrain: |
    An adversary has gained control over a Windows endpoint and has  
    privileges to create scheduled tasks using taskschlr.exe
  targets:
    - Workstations
    - Control Server
    - Input/Output Server
    - Laptop
    - Desktop
    - Remote access
    - Web Application Servers
    - Public-Facing Servers
  platforms:
    - Windows
  severity: Moderate incident
  leverage:
    - Dwelling
    - Infrastructure Compromise
    - Elevation of privilege
  impact:
    - Business disruption
    - Impairement
    - Lose Capabilities
  viability: Likely
  description: |
    A threat actor can use the Task Scheduler to create tasks for nefarious purposes. For example, during 
    the creation of a task in Task Scheduler > Create Task, the newly created task can be configured with 
    specific options in the tab "Action". Through the "Action" settings the threat actor can specify 
    the actions for the task to execute, which often include outbound connections to attacker 
    infrastructure, binary execution or registry entry editing or creation.
 
    Example: 

    Action: "Start a program"
    In Program/script field: C:\ProgramData\<name.exe>

    Tasks are stored in C:\Windows\System32\Tasks\ in XML format.

    They reside in the following registry:

    HKLM\Software\Microsoft\Windows\CurrentVersion\Schedule\TaskCache\Tasks\{GUID}

    “Actions” value stored within the Tasks\{GUID} key points to the command line associated with the task.
    In the registry under HKLM\Software\Microsoft\Windows\CurrentVersion\Schedule\TaskCache\Tasks\{GUID}
    there is a registry key for the task with a key "Actions" related to Tasks\{GUID} id. In "Edit Binary Value"
    of the registy key "Action" there is a reference to the path which can contain malicious executable file,
    for example: C:\Windows\System32\taskschlr.exe or other.

    Similar information is stored within an extensionless XML file, created within C:\Windows\System32\Tasks, 
    where the name of the file matches the name of the task. 
    
    Example for Settings section in XML file with potentially malicious executable file: 

    <Task version="the_version" xmlns="http_link">
     ...
     </Settings>
     <Actions Context="Author">
       <Exec>
        <Command>C\Windows\System32\taskschlr.exe</Command>
       </Exec>
     <Actions>
    </Task>



    

    
